# software-config/supervisor/backend.ini
[program:backend]
directory=/path_to_backend
command=/root/.pyenv/versions/3.9.11/bin/pdm run gunicorn project_name:app_name -c path_to_gunicorn.py
autostart=true
startsecs=30
autorestart=true
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=root
redirect_stderr=true
numprocs=1
stdout_logfile=/path_to_log/gunicorn/pid.log