# software-config/supervisor/celery_worker.ini
[program:celery]
command=/path_to_script/celery.sh
numprocs=1
directory=/path_to_script
autostart=true
autorestart=true
startsecs=30
startretries=3
exitcodes=0,1
stopsignal=QUIT
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/path_to_log/celery/celery.log
stdout_logfile_maxbytes=64MB
stdout_logfile_backups=4
stdout_capture_maxbytes=1MB
stdout_events_enabled=false